# .config/autostart/qtwitter.desktop

[Desktop Entry]
Type=Application
Name=QTwitter
Exec=/usr/bin/qtwitter
Icon=system-run
Comment=
