# .config/autostart/xmonad.desktop
[Desktop Entry]
Type=Application
Name=XMonad
Exec=xmonad --replace
Icon=system-run
Comment=
X-GNOME-Autostart-enabled=true

